# deploy/theta-monitor.service
[Unit]
Description=Theta GPU Monitoring Agent
Documentation=https://github.com/Asomisetty27/theta
After=network.target nvidia-dcgm.service
Wants=nvidia-dcgm.service

[Service]
Type=simple
User=theta
Group=theta
ExecStart=/usr/local/bin/theta monitor --quiet
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=theta

# Config lives at ~/.theta/config.json for the service user
# or override with: Environment=THETA_CONFIG=/path/to/config.json

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target
